# roles/cupie/files/temp.service
[Unit]
Description=Temperature Logger

[Service]
Type=oneshot
ExecStart=/usr/local/bin/log-temp