# roles/cupie/files/temp.timer
[Timer]
OnCalendar=minutely
Persistent=true